# config/docker/supervisor/supervisord-dev.conf
[supervisord]
logfile=syslog
nodaemon=true
loglevel=trace
pidfile = /tmp/supervisord.pid

[program:php-fpm]
command = /usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
priority=10
stopsignal=QUIT

[program:nginx]
command=/usr/sbin/nginx -e stderr -g "daemon off;"
autostart=true
autorestart=true
priority=20
stopsignal=QUIT

[program:vite]
command=npm run dev
directory=/var/www/supervisor-monitor/assets
autostart=true
autorestart=true
priority=20
stopsignal=QUIT
