# config/docker/supervisor/supervisord-dev.conf
[supervisord]
logfile=syslog
nodaemon=true
loglevel=trace
pidfile = /tmp/supervisord.pid

[program:php-fpm]
command = /usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
priority=10
stopsignal=QUIT

[program:nginx]
command=/usr/sbin/nginx -e stderr -g "daemon off;"
autostart=true
autorestart=true
priority=20
stopsignal=QUIT

[program:vite]
command=npm run dev
directory=/var/www/supervisor-monitor/assets
autostart=true
autorestart=true
priority=20
stopsignal=QUIT

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[rpcinterface:twiddler]
supervisor.rpcinterface_factory = supervisor_twiddler.rpcinterface:make_twiddler_rpcinterface

[inet_http_server]
port=*:$SUPERVISOR_RPC_PORT
username=$SUPERVISOR_RPC_USERNAME
password=$SUPERVISOR_RPC_PASSWORD
