# tests/Functional/resources/supervisord.conf
[supervisord]
logfile=syslog
nodaemon=true
loglevel=info
pidfile = /tmp/supervisord.pid

[program:tail-f-void]
command = /usr/bin/tail -f /dev/null
autostart=true
autorestart=true
priority=10
stopsignal=QUIT
numprocs=1
process_name=%(process_num)s

[program:tail-f-void-multiple]
command = /usr/bin/tail -f /dev/null
autostart=true
autorestart=true
priority=10
stopsignal=QUIT
numprocs=3
process_name=%(program_name)s_%(process_num)02d

[program:date-multiple]
command = date
autostart=true
autorestart=true
priority=10
stopsignal=QUIT
numprocs=3
process_name=%(program_name)s_%(process_num)02d

#[program:tail-f-void-multiple-many]
#command = /usr/bin/tail -f /dev/null
#autostart=true
#autorestart=true
#priority=10
#stopsignal=QUIT
#numprocs=300
#process_name=%(program_name)s_%(process_num)03d

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[rpcinterface:twiddler]
supervisor.rpcinterface_factory = supervisor_twiddler.rpcinterface:make_twiddler_rpcinterface

[inet_http_server]
port=*:9555
username=default
password=default

[supervisorctl]
serverurl=http://127.0.0.1:9555
username=default
password=default
